# storerozbot.service
[Unit]
Description=متجر روز — Telegram Bot (@StoreRozbot)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/storerozbot
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
